# systemd/epipyweb-database-rotate.service
[Unit]
Description=Rotate Epipyweb database

[Service]
Type=oneshot
ExecStart=/usr/sbin/epipyweb-database-rotate
StandardError=syslog
